# software/git/cron/github-download-user-repos.crontab
0 */8 * * * . "$HOME/.bashrc" && "github-download-user-repos.sh" "$CLOUD_FOLDER/me/repo/github"
