# software/git/cron/github-download-user-repos.crontab
0 */8 * * * . "$HOME/.bashrc" && "$HOME/.distro-config/bin/github-download-user-repos.sh" "$CLOUD_FOLDER/me/repo/github"
