# software/git/cron/github-download-starred-repos.crontab
0 */2 * * * "github-download-starred-repos.sh" "$CLOUD_FOLDER/archive/repo/3rd-party/github"
